# pgmcp.service
[Unit]
Description=pgmcp - PostgreSQL MCP File Indexer
After=postgresql.service
Wants=postgresql.service

[Service]
Type=notify
ExecStart=/usr/local/bin/pgmcp daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=%i
Environment=PGMCP_CONFIG=%h/.config/pgmcp/config.toml
# Bound glibc thread arenas. Default is 8 × num_cpus (= 512 on a 64-CPU host).
# Each arena retains its high-water mark forever, so a transient burst of
# allocations across many threads bloats RSS by tens of GB that never get
# returned to the kernel. Capping at 2 keeps total resident memory close to
# the actual live working set without measurably hurting throughput.
Environment=MALLOC_ARENA_MAX=2

[Install]
WantedBy=default.target
