# pgmcp.service
[Unit]
Description=pgmcp - PostgreSQL MCP File Indexer
After=postgresql.service
Wants=postgresql.service

[Service]
Type=notify
ExecStart=/usr/local/bin/pgmcp daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
User=%i
Environment=PGMCP_CONFIG=%h/.config/pgmcp/config.toml
Environment=LD_LIBRARY_PATH=/opt/intel/oneapi/mkl/latest/lib/intel64

[Install]
WantedBy=default.target
